# wayland-keyboard-backlight/led-permission.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/local/bin/led-permission.sh

[Install]
WantedBy=default.target